# docker/crontab
* * * * * /speedtest/Speedtest/Speedtest inode | /speedtest/Speedtest.InfluxDb/Speedtest.InfluxDb >> /var/log/cron.log 2>&1

