# docker/crontab
*/15 * * * * /speedtest/speedtest.sh >> /var/log/cron.log 2>&1

